# system/mint/xfce/root/etc/systemd/system/greetd.service.d/override.conf
# this ensures boot logs don't overlap with the greeter, tuigreet
[Service]
Type=idle

