# system/mx/xfce/root/etc/profile.d/nixenv.sh
PATH="/nix/var/nix/profiles/default/bin:$PATH"
export PATH

